# files/usr/lib/systemd/system/netf-bootstrap.service
[Unit]
Description=Piotr first-boot core bootstrap (LUKS enrollment, recovery, firmware)
After=systemd-user-sessions.service network-online.target
Wants=network-online.target
ConditionPathExists=!/var/lib/bootstrap/.core.done

[Service]
Type=oneshot
ExecStart=/usr/share/bootstrap/run-profile.sh core
RemainAfterExit=yes
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
